# supervisor/conf.d/appsec-timeout.conf
[program:appsec-timeout]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/appsec/artisan queue:work database --queue=timeout --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=20
redirect_stderr=true
stdout_logfile=/var/www/html/appsec/queue.log
stopwaitsecs=3600
